# features.in/network-protocols-in-linux/live/files/etc/X11/profile.d/ssh-agent.sh
#!/bin/sh

[ -n "$SSH_AUTH_SOCK" ] ||
	eval `ssh-agent -u 2>/dev/null`
