# features.in/02protocols/files/etc/systemd/system/avahi-daemon.socket.d/override.conf
[Socket]
ExecStartPre=rm -rf /run/avahi-daemon /var/resolv/run/avahi-daemon
ExecStartPre=mkdir -p /var/resolv/run/avahi-daemon
ExecStartPre=ln -sr /var/resolv/run/avahi-daemon /run/avahi-daemon

