# features.in/02protocols/files/etc/profile.d/0lang.sh
# /etc/profile.d/0lang.sh - set i18n stuff
LC_ALL=C
LC_COLLATE=C
LANG=C
